# root/etc/services.d/UnicornTranscoder/run
#!/usr/bin/with-contenv bash

echo "Starting UnicornTranscoder."

exec s6-setuidgid plex /bin/sh -c 'cd /opt/UnicornTranscoder* && npm start'
