# root/etc/services.d/UnicornTranscoder/run
#!/usr/bin/with-contenv bash

echo "Starting UnicornTranscoder."

export INSTANCE_ADDRESS="http://$(hostname -s):3000"

exec s6-setuidgid plex /bin/sh -c 'cd /opt/UnicornTranscoder* && npm start'
